# ha-integration/addon/sentihome/rootfs/etc/services.d/notify/run
#!/usr/bin/with-contenv bash
# SentiHome notify dispatcher (push, TTS, ask flow).
exec /app/.venv/bin/python -m sentihome_notify
