# ha-integration/addon/sentihome/rootfs/etc/services.d/vlm-router/run
#!/usr/bin/with-contenv bash
# SentiHome VLM router (multi-backend inference dispatcher).
exec python -m sentihome_vlm_router
