# ha-integration/addon/sentihome/rootfs/etc/services.d/core/run
#!/usr/bin/with-contenv bash
# Long-running SentiHome core (triage + rule engine + action dispatcher).
exec python -m sentihome_core
