# ha-integration/addon/sentihome/rootfs/etc/services.d/preprocessor/run
#!/usr/bin/with-contenv bash
# SentiHome preprocessor (MOG2 motion + on-camera AI corroboration).
exec python -m sentihome_preprocessor
