# ha-integration/addon/sentihome/rootfs/etc/services.d/ha-agent/run
#!/usr/bin/with-contenv bash
# SentiHome ha-agent (read+write HA bridge).
exec python -m sentihome_ha_agent
